# src/server/etc/cron.d/hmanager-server
0 * * * *  root /usr/bin/perl /usr/local/bin/gen_conf.pl 
