# src/client/etc/cron.d/hmanager
* * * * * root /bin/bash  /usr/local/bin/sync_html.sh
0 * * * *  root /bin/bash  /usr/local/bin/sync_conf.sh
